# server/dev_services/proxy/services/tinyproxy/run
#!/usr/bin/env bash
exec /usr/bin/tinyproxy -d -c /srv/tinyproxy.conf