# sensor/platforms/bbb/manager.service
[Unit]
Description=HoneySens Sensor Manager

[Service]
ExecStart=/opt/manager/venv/bin/manager -p bbb /opt/config.tar.gz

[Install]
WantedBy=multi-user.target
