# server/utils/dev_services/proxy/services/hostapd/run
#!/usr/bin/env bash
exec /usr/sbin/hostapd -i eth0 /srv/hostapd.conf