# server/utils/dev_services/proxy/services/hostapd_listener/run
#!/usr/bin/env bash
exec /usr/sbin/hostapd_cli -a /srv/hostapd_handler.sh